# crontab
* *	* * *	root	/root/auto-claim-rewards.ds
